# deploy/com.quanthaha.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quanthaha.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jason/Documents/Github/QuantHaHa/engine/target/release/server</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jason/Documents/Github/QuantHaHa/engine</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/qhh-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/qhh-server.log</string>
</dict>
</plist>
